# scripts/com.necmttn.ax-checkpoint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.necmttn.ax-checkpoint</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd __AX_DIR__ &amp;&amp; bun src/cli/index.ts ingest --stages=opportunities &amp;&amp; bun src/cli/index.ts improve checkpoint &gt;&gt;__LOG_DIR__/checkpoint.log 2&gt;&amp;1</string>
  </array>
  <!-- Weekly (every 7 days) at 09:05 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/checkpoint.out</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/checkpoint.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__HOME__/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
